# build/sysroot/usr/lib/systemd/system/casaos-message-bus.service
[Unit]
After=casaos-gateway.service
Description=CasaOS Message Bus Service

[Service]
ExecStartPre=/usr/bin/casaos-message-bus -v
ExecStart=/usr/bin/casaos-message-bus -c /etc/casaos/message-bus.conf
PIDFile=/var/run/casaos/message-bus.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
